# build/etc/services.d/mssql/run
#!/usr/bin/with-contenv sh

set -e

ACCEPT_EULA=${ACCEPT_EULA:-"Y"}
MSSQL_SA_PASSWORD=${MSSQL_SA_PASSWORD:-"yourStrong(!)Password"}
MSSQL_PID=${MSSQL_PID:-"Express"}
MSSQL_COLLATION=${MSSQL_COLLATION:-""}

HELPER=""
if [ -f "/opt/mssql/bin/nonroot_msg.sh" ]; then
  # For 2017 MSSQL release
  HELPER="/opt/mssql/bin/nonroot_msg.sh"
elif [ -f /opt/mssql/bin/permissions_check.sh ]; then
  # For 2019 MSSQL release
  HELPER="/opt/mssql/bin/permissions_check.sh"
fi

if [ ! -z ${MSSQL_COLLATION} ] && [ -d /var/opt/mssql/data ]; then
  # change server collation
  s6-setuidgid mssql ${HELPER} /opt/mssql/bin/sqlservr -q${MSSQL_COLLATION}
fi

exec s6-setuidgid mssql ${HELPER} /opt/mssql/bin/sqlservr
