# build/etc/services.d/mssql/run
#!/usr/bin/with-contenv sh

set -e

ACCEPT_EULA=${ACCEPT_EULA:-"Y"}
MSSQL_SA_PASSWORD=${MSSQL_SA_PASSWORD:-"yourStrong(!)Password"}
MSSQL_PID=${MSSQL_PID:-"Express"}
MSSQL_COLLATION=${MSSQL_COLLATION:-""}

if [ ! -z ${MSSQL_COLLATION} ] && [ -d /var/opt/mssql/data ]; then
  # change server collation
  s6-setuidgid mssql /opt/mssql/bin/sqlservr -q${MSSQL_COLLATION}
fi

exec s6-setuidgid mssql /opt/mssql/bin/permissions_check.sh /opt/mssql/bin/sqlservr
